# bin/kittenMQ.service
[Unit]
Description=Kitten-mq
After=network.target

[Service]
User=kitten-mq
ExecStart=@exec start
Restart=always

[Install]
WantedBy=multi-user.target
